# recipes-azure-iot/iot-hub-device-update/iot-hub-device-update/adu-agent.timer
[Unit]
Description=Restart adu-agent

[Timer]
OnUnitInactiveSec=10min

[Install]
WantedBy=timers.target