# recipes-azure-iot/iot-hub-device-update/iot-hub-device-update/deviceupdate-agent.timer
[Unit]
Description=Restart deviceupdate-agent

[Timer]
OnBootSec=5min
OnUnitInactiveSec=10min

[Install]
WantedBy=timers.target
