# recipes-azure-iot/iot-hub-device-update/iot-hub-device-update/deviceupdate-agent.path
[Unit]
Before=paths.target

[Path]
PathChanged=/run/omnect-device-service/omnect_validate_update

[Install]
WantedBy=paths.target
