# recipes-core/systemd/systemd/ics-dm-first-boot.service
[Unit]
Description=first boot
ConditionFirstBoot=yes

[Service]
Type=oneshot
ExecStart=/usr/bin/ics_dm_first_boot.sh

[Install]
WantedBy=multi-user.target
