# recipes-omnect/omnect-base-files/omnect-base-files/etc/profile.d/tpm2tools.sh
export TPM2TOOLS_TCTI="device:/dev/tpmrm0"
